# docker/rootfs/etc/s6-overlay/s6-rc.d/system-config/up
/etc/s6-overlay/scripts/system-config.sh